# config/templates/upstart/ovs-webapp-api.conf
description "webapp api"

start on started ovs-rabbitmq and started ovs-memcached and started ovs-arakoon-ovsdb
stop on stopping ovs-rabbitmq or stopping ovs-memcached or stopping ovs-arakoon-ovsdb

kill timeout 60
respawn
respawn limit 10 5
console log

env PYTHONPATH=:/opt/OpenvStorage

chdir /opt/OpenvStorage/webapps
exec /usr/bin/gunicorn -c /opt/OpenvStorage/config/django_gunicorn_ovs.cfg.py api.wsgi
