# config/templates/upstart/ovs-nginx.conf
description "nginx HTTP Server"

start on started ovs-webapp-api
stop on stopping ovs-webapp-api

kill timeout 5

respawn
exec /usr/sbin/nginx
