# config/templates/systemd/ovs-watcher-config.service
[Unit]
Description=Watcher for the configuration cluster
After=local-fs.target network-online.target

[Service]
Type=simple
SyslogIdentifier=%N
Environment=PYTHONPATH=/opt/OpenvStorage
# The watcher will block until all services are running
ExecStartPre=/usr/bin/python /opt/OpenvStorage/ovs/extensions/generic/watcher.py config wait
# The waiter will block until one of the services is not reachable anymore
ExecStart=/usr/bin/python /opt/OpenvStorage/ovs/extensions/generic/watcher.py config check
Restart=on-failure
TimeoutStopSec=3600

[Install]
WantedBy=multi-user.target
